# etc/cron.d/raid_efi_sync_cron
# Run raid efi sync every 5 minutes
*/5 * * * * root /usr/local/bin/raid_efi_sync
